# systemd/willow.service
[Unit]
Description=Willow Voice Assistant Service
Documentation=https://github.com/saim/willow
After=graphical-session.target

[Service]
Type=dbus
BusName=com.github.saim.Willow
ExecStart=/usr/bin/willow-service
Restart=on-failure
RestartSec=5

# Environment
Environment="PATH=/usr/bin"
Environment="DISPLAY=:0"
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=/run/user/%U"
Environment="DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus"
Environment="XDG_SESSION_TYPE=wayland"

# Security - Relaxed for GUI app launching
PrivateTmp=false
NoNewPrivileges=false

[Install]
WantedBy=default.target
